# dryermon.conf
[program:dryermon]
command=python /home/jdsnape/workers/dryermon/dryermon.py
